# php/7.3-nginx/supervisor/supervisor.d/php-fpm7.conf.default
[program:php-fpm7]
command=/usr/sbin/php-fpm7.3 --nodaemonize --fpm-config=/etc/php/7.3/fpm/pool.d/www.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0